# scripts/server/ai-rss-cleanup.service
[Unit]
Description=AI RSS Artifact Cleanup
After=network-online.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/opt/ai-RSS-person
ExecStart=/bin/bash /opt/ai-RSS-person/scripts/server/cleanup_artifacts.sh
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/ai-RSS-person/reports /opt/ai-RSS-person/logs

[Install]
WantedBy=multi-user.target
