# scripts/server/ai-rss-cookie-rotate.service
[Unit]
Description=AI RSS Twitter Cookie Rotate
After=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/ai-RSS-person
ExecStart=/bin/bash -lc "set -o pipefail; /bin/bash /opt/ai-RSS-person/scripts/server/rotate_twitter_cookie.sh 2>&1 | tee -a /opt/ai-RSS-person/logs/ai-rss-cookie-rotate.log"
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
